# deploy/systemd/vuzol-executor.service
[Unit]
Description=Vuzol isolated coding executor
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=vuzol-executor
Group=vuzol-executor
EnvironmentFile=/etc/vuzol/executor.env
RuntimeDirectory=vuzol/proxy
RuntimeDirectoryMode=0700

# The rootless Docker daemon is deployed as a linger-enabled systemd USER unit
# (see deploy/systemd/user/vuzol-rootless-docker.service).
#
# There is intentionally no After=/Wants= on "vuzol-rootless-docker.service"
# because that name lives in a different systemd manager (the user's).
# A direct cross-manager dependency would be incorrect.
#
# Readiness is fail-closed and enforced in two layers:
# 1. This ExecStartPre waits for the socket path (from env or dynamic user
#    runtime dir) to exist and be a Unix socket. The path is resolved
#    without hard-coding any UID.
# 2. The executor process itself calls RootlessDockerRuntime.preflight()
#    (src/vuzol/execution/sandbox.py) which additionally verifies:
#      - socket is owned by the current euid (vuzol-executor)
#      - daemon reports rootless + seccomp
#      - cgroup v2
#      - no "No memory limit support" / "No cpu cfs quota support" warnings
#    If any check fails the process exits and the unit restarts (on-failure).
#
# This preserves the strict "no rootful, no weak limits" contract.
# ExecStartPre performs a *bounded* wait (max ~30s) for the *socket to exist as Unix socket*.
# It resolves the socket from the configured VUZOL_ env (populated from settings at deploy)
# or dynamically via id -u (portable, no hard-coded UID, runs under User= context).
# It does *not* perform owner check, rootless/seccomp/cgroupv2 validation or limit support
# checks — those are done fail-closed by RootlessDockerRuntime.preflight() inside the
# Python process (see src/vuzol/execution/sandbox.py). If preflight fails the process
# exits and unit restarts.
ExecStartPre=/bin/sh -c "\
  set -e; \
  SOCKET=\"$${VUZOL_EXECUTION__ROOTLESS_DOCKER_SOCKET:-}\"; \
  if [ -z \"$$SOCKET\" ]; then SOCKET=\"/run/user/$$(id -u)/docker.sock\"; fi; \
  echo \"vuzol-executor: waiting for rootless Docker socket at $$SOCKET\"; \
  for i in $$(seq 1 150); do [ -S \"$$SOCKET\" ] && exit 0; sleep 0.2; done; \
  echo 'timeout waiting for rootless docker socket' >&2; exit 1 \
"

ExecStart=/opt/vuzol/.venv/bin/vuzol-executor
TimeoutStartSec=60
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/srv/vuzol/worktrees /srv/vuzol/artifacts /var/lib/vuzol-provider-state
ReadOnlyPaths=/srv/vuzol/repositories /etc/vuzol

[Install]
WantedBy=multi-user.target
