# deploy/systemd/vuzol-applier.service
[Unit]
Description=Vuzol approved-result control and apply worker
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=vuzol-executor
Group=vuzol-executor
EnvironmentFile=/etc/vuzol/executor.env
ExecStart=/opt/vuzol/.venv/bin/vuzol-applier
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/srv/vuzol/repositories
ReadOnlyPaths=/srv/vuzol/worktrees /srv/vuzol/artifacts /etc/vuzol

[Install]
WantedBy=multi-user.target
