# deploy/systemd/vuzol-rootless-docker.service
# LEGACY / MIGRATION REFERENCE ONLY
#
# This was the former system-level unit definition:
#   Type=simple + User=vuzol-executor + RuntimeDirectory + hard-coded XDG path.
#
# PRODUCTION DEPLOYMENTS MUST USE THE USER UNIT:
#   deploy/systemd/user/vuzol-rootless-docker.service
# installed to /etc/systemd/user/vuzol-rootless-docker.service
# and managed via `systemctl --user` (under the vuzol-executor identity with linger).
#
# See README.md (Execution section) and the Step 08 handoff reports for the
# exact installation, readiness, and rollback procedures.
#
# The original body is retained below strictly for reference and rollback.
# Do not install or enable this file as the active daemon in new deployments.
# It is no longer advertised as the production rootless daemon layout.

[Unit]
Description=LEGACY rootless Docker daemon (system unit form - see user/ unit)
After=network.target

[Service]
Type=simple
User=vuzol-executor
Group=vuzol-executor
Environment=HOME=/var/lib/vuzol-executor
Environment=XDG_RUNTIME_DIR=/run/vuzol-executor
RuntimeDirectory=vuzol-executor
RuntimeDirectoryMode=0700
ExecStart=/usr/bin/dockerd-rootless.sh --host=unix:///run/vuzol-executor/docker.sock
ExecStartPost=/bin/sh -c 'until /usr/bin/docker --host unix:///run/vuzol-executor/docker.sock info >/dev/null 2>&1; do sleep 0.2; done'
Restart=always
RestartSec=2
Delegate=yes
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
