# deploy/systemd/user/vuzol-rootless-docker.service
[Unit]
Description=Dedicated rootless Docker daemon for Vuzol sandboxes
Documentation=man:dockerd-rootless(1)
After=network.target

[Service]
Type=simple
# USER UNIT: runs under the linger-enabled vuzol-executor user's systemd manager.
# Do NOT include User= or Group= (user units inherit the user identity).
# The user manager provides XDG_RUNTIME_DIR (exposed as %t) and the delegated
# user slice for cgroup controllers (cpu, memory, pids, etc.).
# This enables reliable enforcement of --cpus, --memory, --pids-limit etc.
# Socket path uses %t for portability — never hard-code a numeric UID here.
ExecStart=/usr/bin/dockerd-rootless.sh --host=unix://%t/docker.sock
ExecStartPost=/bin/sh -c 'until /usr/bin/docker --host unix://%t/docker.sock info >/dev/null 2>&1; do sleep 0.2; done'
Restart=always
RestartSec=2
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=default.target
